# src/daemon/systemd/user/pulseaudio.service.in
[Unit]
Description=Sound Service

[Service]
Environment=JACK_PROMISCUOUS_SERVER=audio
ExecStart=@PA_BINARY@ --daemonize=no
Restart=on-failure

[Install]
Also=pulseaudio.socket
WantedBy=default.target
